# api.service
[Unit]
Description=API server

[Service]
ExecStart=/setup.sh

[Install]
Alias=api.service
